# templates/galaxy.fedora-init.erb
#!/bin/bash
#
# Init file for Galaxy (http://galaxyproject.org/)
#   Suitable for use on Fedora and derivatives (RedHat Enterprise Linux, Scientific Linux, CentOS)
#
# Contributed by Brad Chapman
#
# chkconfig: 2345 98 20
# description: Galaxy http://galaxyproject.org/

#--- config

SERVICE_NAME="Galaxy"
RUN_AS="galaxy"
RUN_IN="<%= @directory %>"

#--- main actions

start() {
	echo -e "Starting $SERVICE_NAME...\n"
	 <% if @wk_config %>
  	su - $RUN_AS -c 'GALAXY_RUN_ALL=42 <%= @directory %>/run.sh --daemon'
        <% else %>
        su - $RUN_AS -c '<%= @directory %>/run.sh --daemon'
        <% end %>	
        echo "...done."
}

stop() {
	echo -ne "Stopping $SERVICE_NAME...\n"
         <% if @wk_config %>
  	su - $RUN_AS -c 'GALAXY_RUN_ALL=42 <%= @directory %>/run.sh --stop-daemon'
        <% else %>
        su - $RUN_AS -c '<%= @directory %>/run.sh --stop-daemon'
        <% end %>	
        echo "done."
}

status() {
	echo -n "$SERVICE_NAME status: "
	
	while read pid; do
		if [ "$(readlink -m /proc/$pid/cwd)" = "$(readlink -m $RUN_IN)" ]; then
			echo "started"
			return 0
		fi
	done < <(ps ax -o 'pid cmd' | grep -P '^\s*\d+ python ./scripts/paster.py serve' | awk '{print $1}')
	echo "stopped"
	return 3
}

notsupported() {
	echo "*** ERROR*** $SERVICE_NAME: operation [$1] not supported"
}

usage() {
	echo "Usage: $SERVICE_NAME start|stop|restart|status"
}


#---

case "$1" in
	start)
		start "$@"
		;;
	stop)
		stop
		;;
	restart|reload)
		stop
		start
		;;
	status)
		set +e
		status
		exit $?
		;;
	'')
		usage >&2
		exit 1
		;;
	*)
		notsupported "$1" >&2
		usage >&2
		exit 1
		;;
esac
